# root/etc/s6-overlay/s6-rc.d/init-restic-test/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -e 

restic -r "$RESTIC_REPOSITORY" snapshots --host="$RESTIC_HOSTNAME"
