# root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec \
    s6-notifyoncheck -d -n 300 -w 1000 \
    /usr/sbin/crond -f -S -l 5 -c /config/crontabs
